# supervisord.conf
[supervisord]
nodaemon=true

[program:nxfilter]
command=/bin/bash -c "/nxfilter/bin/startup.sh"

[program:cron]
command=/bin/bash -c "/usr/sbin/cron -f"